# deploy/systemd/wa-watchdog.timer
[Unit]
Description=Run wa-watchdog every 2 minutes (recover hung instance containers)

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
AccuracySec=15s
Unit=wa-watchdog.service

[Install]
WantedBy=timers.target
